# toolchain/btx/btx-miner-guard.service
[Unit]
Description=Pause BTX mining during interactive/build work
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/btx-miner-guard
